# ops/com.drivers-school.nightly-readiness.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>com.drivers-school.nightly-readiness</string>
  <key>WorkingDirectory</key><string>/Users/clpc/drivers-school</string>
  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/npx</string>
    <string>tsx</string>
    <string>--conditions=react-server</string>
    <string>scripts/nightly-readiness.ts</string>
  </array>
  <key>EnvironmentVariables</key>
  <dict><key>PATH</key><string>/opt/homebrew/bin:/usr/bin:/bin:/usr/sbin:/sbin</string></dict>
  <key>StartCalendarInterval</key>
  <dict><key>Hour</key><integer>3</integer><key>Minute</key><integer>30</integer></dict>
  <key>StandardOutPath</key><string>/Users/clpc/drivers-school/var/log/nightly-readiness.log</string>
  <key>StandardErrorPath</key><string>/Users/clpc/drivers-school/var/log/nightly-readiness.log</string>
</dict>
</plist>
